# task/etl.cron

01 01 * * * . /home/hadoop/ihealth/task/env.sh && /bin/bash /home/hadoop/ihealth/task/import-data.sh 120.131.8.76 ih_aijiakang root Wifah0Aaj3 ih_blood_pressure >> /home/hadoop/ihealth/logs/import_data_bp_`date +\%Y-\%m`.log 2>&1
01 02 * * * . /home/hadoop/ihealth/task/env.sh && /bin/bash /home/hadoop/ihealth/task/import-data.sh 120.131.8.76 ih_aijiakang root Wifah0Aaj3 ih_blood_pressure_all >> /home/hadoop/ihealth/logs/import_data_bp_all_`date +\%Y-\%m`.log 2>&1
01 03 * * 1 . /home/hadoop/ihealth/task/env.sh && /bin/bash /home/hadoop/ihealth/task/import-data.sh 120.131.8.76 ih_aijiakang root Wifah0Aaj3 ih_user_profile >> /home/hadoop/ihealth/logs/import_data_user_`date +\%Y-\%m`.log 2>&1
01 04 * * * . /home/hadoop/ihealth/task/env.sh && /bin/bash /home/hadoop/ihealth/task/hdfs-mysql.sh >> /home/hadoop/ihealth/logs/migrate_to_mysql_`date +\%Y-\%m`.log 2>&1
1 5 * * * . /home/hadoop/ihealth/task/env.sh && /bin/bash /home/hadoop/ihealth/task/loadData.sh 10.9.0.11:32400 ih_aijiakang.ih_blood_pressure >> /home/hadoop/ihealth/logs/load_data_bp_`date +\%Y-\%m`.log 2>&1
1 5 * * * . /home/hadoop/ihealth/task/env.sh && /bin/bash /home/hadoop/ihealth/task/loadData.sh 10.9.0.11:32400 ih_aijiakang.ih_blood_pressure_all >> /home/hadoop/ihealth/logs/load_data_bp_all_`date +\%Y-\%m`.log 2>&1
1 5 * * 1 . /home/hadoop/ihealth/task/env.sh && /bin/bash /home/hadoop/ihealth/task/loadData.sh 10.9.0.11:32400 ih_aijiakang.ih_user_profile >> /home/hadoop/ihealth/logs/load_data_user_`date +\%Y-\%m`.log 2>&1
1 5 * * * . /home/hadoop/ihealth/task/env.sh && /bin/bash /home/hadoop/ihealth/task/alter-ptt.sh >> /home/hadoop/ihealth/logs/alter_ptt_`date +\%Y-\%m`.log 2>&1
1 5 * * * . /home/hadoop/ihealth/task/env.sh && /bin/bash /home/hadoop/ihealth/task/insert-to-ptt.sh >> /home/hadoop/ihealth/logs/insert_to_ptt_`date +\%Y-\%m`.log 2>&1
1 6 * * * . /home/hadoop/ihealth/task/env.sh && /bin/bash /home/hadoop/ihealth/task/weeklyTask.sh /home/hadoop/ihealth/aijiakang_data_analysis.jar >> /home/hadoop/ihealth/logs/weekly_task_`date +\%Y-\%m`.log 2>&1

