# task/weekly-report.cron

36 15 * * * . /home/hadoop/ihealth/task/env.sh && /bin/bash /home/hadoop/ihealth/task/weekly-report.sh /home/hadoop/ihealth/aijiakang_data_analysis.jar >> /home/hadoop/ihealth/logs/weekly_task_`date +\%Y-\%m`.log 2>&1
